# supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/etc/init.d/apache2 start
autostart=true
autorestart=true